# profile.d/40-DE
# Set locale for legacy desktop environments that do not respect localectl.
# TODO: export and fallback all the stuffs as defined in locale.conf(7).
. /etc/locale.conf
: ${LC_ALL=${LANG=C}}
export LANG LC_ALL
if test -z "${XDG_RUNTIME_DIR}"; then
    export XDG_RUNTIME_DIR=/tmp/${UID}-runtime-dir
    if ! test -d "${XDG_RUNTIME_DIR}"; then
        mkdir "${XDG_RUNTIME_DIR}"
        chmod 0700 "${XDG_RUNTIME_DIR}"
    fi
fi

