# profile.d/extrapaths
if [ -d ~/bin ]; then
        pathprepend ~/bin
fi
