# profile.d/extrapaths
if [ -d ~/bin ]; then
    PATH="$HOME/bin:$PATH"
fi
